# websitewatcher.timer
[Unit]
Description=websitewatcher timer

[Timer]
OnBootSec=15min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
